# docker/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=sh -c '. /etc/apache2/envvars && apache2ctl -DFOREGROUND'

[program:mysql]
command=mysqld_safe
